# firmware-addons/rc.local
# rc.local, runs on boot
chmod +x /etc/4re5.sh
sh /etc/4re5.sh
